# oneshoot-service.timer
[Unit]
Description={{DESCRIPTION}}

[Timer]
OnCalendar={{ONCALENDAR}}
Unit={{UNIT}}

[Install]
WantedBy={{WANTEDBY}}
