# oneshoot-service.timer
[Unit]
Description={{DESCRIPTION}}

[Timer]
OnBootSec={{ONBOOTSEC}}
OnCalendar={{ONCALENDAR}}

[Install]
WantedBy=default.target
